# crontabs
0 * * * * ~/current/source/manage.py popularity2rank && supervisorctl restart uwsgi
