# crontabs
MAILTO="s.angel+presidapps@twidi.com"
SHELL=/bin/bash
0 * * * * source ~/env/bin/activate && ~/current/source/manage.py popularity2rank && supervisorctl restart uwsgi
